# btrfs-snapshot/system-sec\x2dcounter.slice
[Unit]
Description=Slice for sec-counter@.service template
SurviveFinalKillSignal=yes
IgnoreOnIsolate=yes
DefaultDependencies=no
